# .docker/rootfs/etc/supervisor.d/supervisord.ini
[supervisord]
nodaemon = true

[program:nginx]
command=nginx

[program:php-fpm]
command=php-fpm
stdout_logfile=/var/log/php.info.log
stderr_logfile=/var/log/php.error.log
stdout_logfile_maxbytes=16MB
stderr_logfile_maxbytes=16MB
